# edge/quadlet/openclaw-agent.container
# ~/.config/containers/systemd/openclaw-agent.container
#
# OpenClaw edge agent — runs as a rootless podman Quadlet managed by systemd.
# Stopped by default (Restart=no). Central supervisor starts/stops via SSH.
#
# Install: cp openclaw-agent.container openclaw-config.volume ~/.config/containers/systemd/
#          systemctl --user daemon-reload
# Start:   systemctl --user start openclaw-agent
# Logs:    journalctl --user -u openclaw-agent -f

[Unit]
Description=OpenClaw Edge Agent

[Container]
Image=quay.io/sallyom/openclaw:latest
ContainerName=openclaw-agent

# Config + workspace persist across activations
Volume=openclaw-config.volume:/home/node/.openclaw:Z

# Host networking so agent can reach local services, databases, APIs
Network=host

# Runtime environment
Environment=NODE_ENV=production
Environment=OPENCLAW_GATEWAY_TOKEN=${OPENCLAW_GATEWAY_TOKEN}

# Gateway entrypoint
Exec=node openclaw.mjs gateway --allow-unconfigured

[Service]
# CRITICAL: agent stays stopped until supervisor explicitly starts it.
# No automatic restarts — central gateway controls lifecycle.
Restart=no
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=default.target
